# crontab
#Put this line into /etc/crontab
*/15      *      *       *       *       postgres    psql -c "select delete_old_rows()" knowhub_test